# deploy/bq-lora-ui-shutdown.service
[Unit]
Description=Splash d'arret/redemarrage sur l'ecran SPI
# IMPORTANT : on garde les dependances par defaut pour que systemd ajoute
# Before=shutdown.target + Conflicts=shutdown.target. Sans elles, le service
# n'est pas ordonne dans la sequence d'arret et son ExecStop ne s'affiche pas.
# Before=bq-lora-ui.service -> a l'arret, ce service est stoppe APRES bq-lora-ui,
# donc son ExecStop dessine une fois que bq-lora-ui a libere /dev/fb0.
Before=bq-lora-ui.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/usr/local/sbin/bq-lora-ui-shutdown-splash
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
