# deploy/bq-lora-ui-splash.service
[Unit]
Description=Boot splash sur l'ecran SPI (avant bq-lora-ui)
DefaultDependencies=no
After=local-fs.target
Before=bq-lora-ui.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Garantit /dev/fb_spi -> panneau SPI (par NOM) et attend son apparition.
# Remplace l'ancien 'until [ -e /dev/fb0 ]' qui bloquait pour toujours quand le
# panneau atterrit sur fb1 (course simplefb/fbtft) -> ecran noir intermittent.
ExecStartPre=/usr/local/sbin/bq-lora-ui-fb-spi-link
# Detache la console tot -> pas de texte de boot sur l'ecran
ExecStartPre=/bin/sh -c 'echo 0 > /sys/class/vtconsole/vtcon1/bind || true'
ExecStart=/usr/bin/python3 /home/bq-lora/bq-lora-ui/tools/splash.py

[Install]
WantedBy=sysinit.target
