# deploy/meshui-shutdown.service
[Unit]
Description=Splash d'arret/redemarrage sur l'ecran SPI
# IMPORTANT : on garde les dependances par defaut pour que systemd ajoute
# Before=shutdown.target + Conflicts=shutdown.target. Sans elles, le service
# n'est pas ordonne dans la sequence d'arret et son ExecStop ne s'affiche pas.
# Before=meshui.service -> a l'arret, ce service est stoppe APRES meshui,
# donc son ExecStop dessine une fois que meshui a libere /dev/fb0.
Before=meshui.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/usr/local/sbin/meshui-shutdown-splash
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
